# linux/.config/runit/sv/ssh-agent/run
#!/bin/sh
rm -f "$HOME/.ssh/ssh_auth_sock"
exec ssh-agent -D -a "$HOME/.ssh/ssh_auth_sock"
