# .config/runit/sv/ssh-agent/run
#!/bin/sh
exec ssh-agent -D -a "$HOME/.ssh/ssh_auth_sock"
